# scripts/services/gh-monitor.service
[Unit]
Description=GitHub monitor for @claude PR comment automation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/puma/Repos/claude-dot-files/scripts/services/gh-monitor.sh
Environment=HOME=/home/%u

[Install]
WantedBy=default.target
