# scripts/services/gh-monitor.service
[Unit]
Description=GitHub monitor for @claude PR comment automation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Path is written by install.sh --with-services at install time
# Do not edit manually — re-run install.sh if the repo moves
ExecStart=INSTALL_PLACEHOLDER
Environment=HOME=/home/%u

[Install]
WantedBy=default.target
